# udev-rules/00-applekb-compose.rules
KERNEL=="hidraw*", SUBSYSTEM=="hidraw", ATTRS{address}=="00:22:41:d9:06:bd", RUN+="setxkbmap -option compose:ralt"